# settings/.config/autostart/lightscreen.desktop
[Desktop Entry]
Exec=/usr/bin/lightscreen -h
Type=Application